# deploy/roles/pg-backup/templates/pg-backup.service.j2
[Unit]
Description=PostgreSQL Database Backup
Documentation=https://github.com/jcfischer/pg-backup
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ pg_backup_user }}
Group={{ pg_backup_group }}

# Load environment from config file
EnvironmentFile=-/etc/pg-backup/pg-backup.env

# Run backup
ExecStart={{ pg_backup_bin_path }} backup

# Allow up to 2 hours for large backups
TimeoutStartSec=7200

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths={{ pg_backup_dir }}
{% if pg_backup_directories %}
{% for dir in pg_backup_directories.split(',') %}
ReadOnlyPaths={{ dir | trim }}
{% endfor %}
{% endif %}

[Install]
WantedBy=multi-user.target
