# deploy/roles/pg-backup/templates/pg-backup.service.j2
[Unit]
Description=PostgreSQL Database Backup
Documentation=https://github.com/jcfischer/pg-backup
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ pg_backup_user }}
Group={{ pg_backup_group }}

# Load environment from config file
EnvironmentFile=/etc/pg-backup/pg-backup.env

# Run backup
ExecStart={{ pg_backup_bin_path }} backup

# Allow up to 2 hours for large backups
TimeoutStartSec=7200

# Minimal security - directory backups need full /home access via ACL
ReadWritePaths={{ pg_backup_dir }}

[Install]
WantedBy=multi-user.target
